# src/etc/rc.local
sudo -u lounge /lounge/triggers/sound reload
exit 0
