# src/etc/rc.local
sudo -u lounge /lounge/triggers/sound reload
echo performance > /sys/devices/system/cpu/cpu0/cpufreq/scaling_governor

echo 1008000 > /sys/devices/system/cpu/cpu0/cpufreq/scaling_max_freq
echo 1008000 > /sys/devices/system/cpu/cpu0/cpufreq/scaling_min_freq

echo on > /sys/devices/platform/mali_dev.0/power/control
echo on > /sys/devices/platform/mali_drm.0/power/control
echo on > /sys/devices/platform/mali_drm.1/power/control
echo on > /sys/class/graphics/fb0/power/control
echo on > /sys/class/graphics/fb1/power/control
echo on > /sys/class/graphics/fbcon/power/control
echo on > /sys/class/sound/card0/power/control
echo on > /sys/class/sound/card1/power/control

touch /var/run/player/playlist_cache
touch /var/run/player/playlist_last
chown -R lounge:lounge /var/run/player/

# disable fb0 screenblanking
echo -e '\033[9;0]' > /dev/tty1
echo "RC local done"
exit 0

